# archives/nodejs-chat-system-complete/roflfaucet-chat.service
[Unit]
Description=ROFL Faucet Chat Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/roflfaucet
ExecStart=/usr/bin/node chat-server-advanced.js
Restart=always
RestartSec=3
StartLimitBurst=5
StartLimitInterval=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=roflfaucet-chat

# Environment variables
Environment=NODE_ENV=production
Environment=NODE_OPTIONS="--max-old-space-size=256"

# Improved shutdown handling
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
TimeoutStopSec=12
TimeoutStartSec=30

# Security and resource limits
NoNewPrivileges=yes
LimitNOFILE=65536
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
